# appfactory-bot.service
[Unit]
Description=AppFactory Telegram Bot
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/appfactory-bot
ExecStart=/usr/bin/python3 -m bot.main
Restart=always
RestartSec=10
EnvironmentFile=/opt/appfactory-bot/.env

# Security hardening
NoNewPrivileges=false
ProtectHome=false

[Install]
WantedBy=multi-user.target
